# virtualMachines/205-SmokePing/public-ip-logging/crontab
0 0 * * * echo $(date +\%Y-\%m-\%d) $(curl -4 https://ifconfig.co 2>/dev/null) >> ~/ipv4_public.txt 2>&1
